# deploy/vya-kitchen.service
[Unit]
Description=Vya's Kitchen Flask app
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/vya
EnvironmentFile=/home/ec2-user/vya/.env
ExecStart=/usr/bin/python3 /home/ec2-user/vya/server.py
Restart=on-failure
RestartSec=3
StandardOutput=append:/home/ec2-user/vya/server.log
StandardError=append:/home/ec2-user/vya/server.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
